# systemd/show-ai-usage.timer
[Unit]
Description=Show AI Usage — Periodic poll timer (every 5 min)
Documentation=https://github.com/your/project

[Timer]
# First run 2 minutes after boot
OnBootSec=2min
# Subsequent runs 10 minutes after the previous run completes
OnUnitActiveSec=3min
# Catch up after sleep / hibernation
Persistent=true

[Install]
WantedBy=timers.target
