# systemd/show-ai-usage.service
[Unit]
Description=Show AI Usage — Poll subscription usage data
Documentation=https://github.com/your/project
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=@@UV_PATH@@ run --directory @@PROJECT_DIR@@ python -m poller.main --oneshot
WorkingDirectory=@@PROJECT_DIR@@
User=%u
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal
